# udock/supervisor/conf.d/ajenti.conf
[program:ajenti]
command=/usr/bin/python /usr/bin/ajenti-panel -d
stdout_events_enabled=true
stderr_events_enabled=true
autostart=true
autorestart=true