# udock/supervisor/conf.d/cron.conf
[program:cron]
command=cron -f
stdout_events_enabled=true
stderr_events_enabled=true
autostart=true
autorestart=true
