# udock/supervisor/conf.d/logrotate.conf
[program:logrotate]
process_name=rsyslog
command=/bin/bash -c "source /etc/default/rsyslog && source /lib/lsb/init-functions && exec /usr/sbin/rsyslogd"
autostart=true
autorestart=true
