# udock/supervisor/conf.d/memcached.conf
[program:memcached]
command=/usr/bin/memcached -p 11211 -u memcache -m 64 -c 1024
stdout_events_enabled=true
stderr_events_enabled=true
autostart=false
autorestart=unexpected