# udock/supervisor/conf.d/mysqld.conf
[program:mysqld]
command=/usr/sbin/mysqld --datadir=/data/mysql
stdout_events_enabled=true
stderr_events_enabled=true
