# udock/supervisor/conf.d/nginx.conf
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;" -c "/data/nginx/nginx.conf"
stdout_events_enabled=true
stderr_events_enabled=true
autostart=true
autorestart=false