# udock/supervisor/conf.d/nginx.conf
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;" -c "/data/nginx/main-conf/nginx.conf"
autostart=true
autorestart=true
