# udock/supervisor/conf.d/php.conf
[program:php5.6-fpm]
command=/usr/sbin/php-fpm5.6 --nodaemonize
autostart=true
autorestart=true

[program:php7.0-fpm]
command=/usr/sbin/php-fpm7.0 --nodaemonize
autostart=true
autorestart=true
