# udock/supervisor/conf.d/redis.conf
[program:redis]
command=/usr/local/bin/redis-server /data/redis/redis.conf
stdout_events_enabled=true
stderr_events_enabled=true
autostart=true
autorestart=unexpected